# blctl.service
[Unit]
After=network.target

[Service]
ExecStart=/root/blctl

[Install]
WantedBy=multi-user.target
