# blctl_chmod.service
[Unit]
Description=
After=

[Service]
ExecStart=/root/blctl_chmod.sh

[Install]
WantedBy=multi-user.target
