# trash/MBP_special/kbctl.service
[Unit]
After=network.target

[Service]
ExecStart=/root/kbctl.sh

[Install]
WantedBy=multi-user.target
